# deployment/crontab
# m h  dom mon dow   command
MAILTO=represent@opennorth.ca
0 4 * * * /home/represent/represent-env/bin/python /home/represent/app/manage.py updaterepresentatives